# mistral/files/mistral-engine.service
[Unit]
Description=Openstack Workflow Engine Service

[Service]
Type=simple
User=mistral
Group=mistral

ExecStart=/usr/bin/mistral-server --server engine --config-file /etc/mistral.conf

[Install]
WantedBy=multi-user.target
