# mistral/files/mistral-executor.service
[Unit]
Description=Openstack Workflow Executor Service

[Service]
Type=simple
User=mistral
Group=mistral

ExecStart=/usr/bin/mistral-server --server executor --config-file /etc/mistral.conf

[Install]
WantedBy=multi-user.target
